# disk-agent/disk-agent.service
[Unit]
Description=Disk Agent - Binary Server Disk Usage Reporter
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/disk-agent/disk_agent.py --path /data/binaries --port 9090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
